# actual/.config/systemd/user/battery_notify.service
[Unit]
Description=Battery level checker


[Service]
ExecStart=/usr/bin/bash %h/.config/dunst/scripts/battery.sh

[Install]
WantedBy=timers.target
